# firmware/ra-update@ttyLP2.service.d/ConditionFirmware.conf
[Unit]
ConditionFirmware=|device-tree-compatible(chargebyte,imx93-charge-som)
ConditionFirmware=|device-tree-compatible(chargebyte,imx93-charge-control-y)
