# firmware/ra-update@ttyLP1.service.d/Environment.conf
[Service]
Environment="SAFETY_MCU_GPIOCHIP=/dev/gpiochip1"
Environment="SAFETY_MCU_RESET_GPIO=nSAFETY2_RESET_INT"
Environment="SAFETY_MCU_MD_GPIO=SAFETY2_BOOTMODE_SET"
